# simplewall_rules.target
[Unit]
Description=Custom SimpleWall Target
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
